# .config/autostart/autokey.desktop
[Desktop Entry]
Type=Application
Name=Autokey
Comment=Run autokey
Exec=autokey-gtk &
OnlyShowIn=GNOME;
X-GNOME-Autostart-Phase=Application
